# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:crond]
command=/usr/sbin/crond

[program:httpd]
command=/usr/sbin/httpd -k start
